# node-image/systemd/sw-agent.service
[Unit]
Description=SecondWind node agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
StateDirectory=secondwind
Environment=SECONDWIND_AGENT_STATE_FILE=%S/secondwind/agent-state.json
Environment=SECONDWIND_AGENT_CERTIFICATE_FILE=%S/secondwind/node-cert.pem
Environment=SECONDWIND_AGENT_PRIVATE_KEY_FILE=%S/secondwind/node-key.pem
EnvironmentFile=-/etc/secondwind/sw-agent.env
ExecStart=/usr/bin/sw-agent
Restart=on-failure

[Install]
WantedBy=multi-user.target
